# scripts/worker.conf
[program:worker]
process_name=%(process_num)s
command=php /var/www/backend/artisan queue:work --sleep=3 --tries=3 --daemon
user=www
autostart=true
autorestart=true
numprocs=2
stdout_logfile=/dev/worker.stdout.log
stderr_logfile=/dev/worker.stderr.log
startsecs=5
startretries=3
redirect_stderr=true


[supervisord]
